# ursine-explorer.service
[Unit]
Description=Ursine Explorer - ADS-B Aircraft Monitor
After=network.target dump1090-mutability.service
Wants=dump1090-mutability.service

[Service]
Type=simple
User=ursine
Group=ursine
WorkingDirectory=/opt/ursine-explorer
ExecStart=/usr/bin/python3 /opt/ursine-explorer/monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ursine-explorer

[Install]
WantedBy=multi-user.target